# deploy/oled-dashboard.service
[Unit]
Description=OLED status dashboard for the Raspberry Pi bot box
After=local-fs.target network-online.target
Wants=network-online.target

[Service]
Type=simple
User=rpi-bot
Group=rpi-bot
SupplementaryGroups=i2c
EnvironmentFile=-/etc/default/oled-dashboard
ExecStart=/opt/rpi-bot/bin/oled-dashboard
Restart=always
RestartSec=2
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
